# systemd_user/qat.timer
[Unit]
Description=Run quality assurance tool

[Timer]
OnBootSec=2min
Persistent=True

[Install]
WantedBy=timers.target